# server/production/supervisor/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
redirect_stderr=true
stdout_logfile=/var/log/nginx/std.log
autostart=true
autorestart=true
